# hosts/macmini/launchagents/com.francisco.arr-safety-guard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.francisco.arr-safety-guard</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/python3</string>
      <string>/Users/franciscomazzucofilho/Scripts/arr-safety-guard.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>StandardOutPath</key>
    <string>/Users/franciscomazzucofilho/Library/Logs/arr-safety-guard.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/franciscomazzucofilho/Library/Logs/arr-safety-guard.log</string>
  </dict>
</plist>
