# hosts/macmini/launchagents/com.francisco.media-compress.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.francisco.media-compress</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/franciscomazzucofilho/Scripts/media-compress-worker.py</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>900</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/franciscomazzucofilho/Tools/lima/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/franciscomazzucofilho/Library/Logs/media-compress.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/franciscomazzucofilho/Library/Logs/media-compress.err.log</string>
</dict>
</plist>
