# services/datalake/crontab
# Datalake ingest schedules (times in UTC, BRT = UTC-3)
# Nubank: 06:00 BRT = 09:00 UTC
0 9 * * *  PLUGGY_CONFIG_PATH=/data/nubank-config.json /app/run-job.sh nubank "python3 /app/scripts/pluggy-ingest.py"

# ORA 03: 06:05 BRT = 09:05 UTC
5 9 * * *  /app/run-job.sh ora03 "python3 /app/scripts/ora03-ingest.py"

# Jellyfin: 09:00 BRT = 12:00 UTC
0 12 * * *  /app/run-job.sh jellyfin "python3 /app/scripts/jellyfin-ingest.py"
